# packer/template/supervisor/plone5.conf
[program:plone5]
user=vagrant
directory=/home/vagrant/Plone/
command=/home/vagrant/Plone/zinstance/bin/plonectl fg
